# dockers/mongoset.supervisor.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:main]
command = mongod --port 27020 --dbpath /usr/local/app/main --replSet dbrs
directory = /usr/local/app/main
autostart = true
autorestart = true

[program:repl]
command = mongod --port 27021 --dbpath /usr/local/app/repl --replSet dbrs
directory = /usr/local/app/repl
autostart = true
autorestart = true
